# ruiji.ini
[program:ruijinet]

command=dotnet /usr/ruiji.net/RuiJi.Net.Cmd.dll --server.urls=http://*:80
directory=/usr/ruiji.net/
startsecs=30
autostart=true
autorestart=true

environment=ASPNETCORE_ENVIRONMENT=Production

redirect_stderr=false
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/usr/ruiji.net/app.log
stopsignal=INT
user=root

[supervisord]
logfile = /tmp/ruiji.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/ruiji.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = ruiji
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false